# debian/mininoip.cron.hourly
#!/bin/sh
[ -x /usr/sbin/mininoip ] && /usr/sbin/mininoip

