# debian/noipmini.cron.hourly
#!/bin/sh
[ -x /usr/sbin/noipmini ] && /usr/sbin/noipmini

